# installer/systemd/control-surface.service
[Unit]
Description=Control Surface Platform
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/control-surface
Environment=NODE_ENV=production
EnvironmentFile=-/etc/control-surface/config.env
WorkingDirectory=/var/lib/control-surface
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
