# installer/systemd/tib-builder.service
[Unit]
Description=TIB Builder Platform
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tib-builder
EnvironmentFile=-/etc/tib-builder/config.env
WorkingDirectory=/var/lib/tib-builder
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target